# scripts/hetzner/daskuechenhaus-mail-importer.service
[Unit]
Description=Daskuechenhaus Mail Importer (IMAP -> PostgreSQL)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=tenant_daskuechenhaus_app
Group=tenant_daskuechenhaus_app
WorkingDirectory=/opt/daskuechenhaus/mail-importer
Environment=DKH_ADMIN_DATABASE=tenant_daskuechenhaus
Environment=DKH_MAIL_ENV_FILE=/etc/daskuechenhaus/mail.env
Environment=DKH_ADMIN_UPLOAD_ROOT=/var/lib/daskuechenhaus/uploads
ExecStart=/usr/bin/python3 /opt/daskuechenhaus/mail-importer/app.py --once
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/daskuechenhaus
ReadWritePaths=/var/lib/daskuechenhaus/uploads
