# scripts/hetzner/daskuechenhaus-mail-importer.timer
[Unit]
Description=Run the Daskuechenhaus Mail Importer every few minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=3min
AccuracySec=30s
Persistent=true
Unit=daskuechenhaus-mail-importer.service

[Install]
WantedBy=timers.target
